# deploy/paper-trader.service
[Unit]
Description=Paper Trader — Multi-Agent Trading System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blaine
WorkingDirectory=/home/blaine/paper-trader
EnvironmentFile=/home/blaine/paper-trader/.env
ExecStart=/home/blaine/paper-trader/venv/bin/python orchestrator.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/blaine/paper-trader/logs/service.log
StandardError=append:/home/blaine/paper-trader/logs/service.log

[Install]
WantedBy=multi-user.target
